# pi/scripts/pyqt5-gui.service
[Unit]
Description=PyQt5 GUI Application
After=graphical.target

[Service]
Type=simple
User=flex
Group=flex
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/flex/.Xauthority"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
WorkingDirectory=/home/flex/uis/src
ExecStart=/usr/bin/python3 /home/flex/uis/src/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
